# supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command = /usr/local/bin/uwsgi --ini /sites/django/uwsgi.ini --touch-reload /sites/django/reload

[program:nginx]
command = /usr/sbin/nginx
